# btc-karotter-bot.service
[Unit]
Description=BTC hourly Karotter bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/btc-karotter-bot
EnvironmentFile=/opt/btc-karotter-bot/.env
ExecStart=/usr/bin/python3 /opt/btc-karotter-bot/main.py
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
